# templates/upstart.conf.erb
description "Serf agent"

start on runlevel [2345]
stop on runlevel [!2345]

exec <%=@install_path%>serf agent \
    -config-file=<%=@config_file%> >> <%=@log_file%> 2>&1

